# systemd/cosmos-notifications.crontab
0 1 * * * SERVICE_DIR=/home/ummc-user/cosmos/ummc-notifications/ping.py API_BASE_URL="https://api.ummc.cosmos.care/v1" BOT_TOKEN= API_AUTHORIZATION_TOKEN= $SERVICE_DIR/python $SERVICE_DIR/ping.py --session 0
0 4 * * * SERVICE_DIR=/home/ummc-user/cosmos/ummc-notifications/ping.py API_BASE_URL="https://api.ummc.cosmos.care/v1" BOT_TOKEN= API_AUTHORIZATION_TOKEN= $SERVICE_DIR/python $SERVICE_DIR/ping.py --session 1
0 5 * * * SERVICE_DIR=/home/ummc-user/cosmos/ummc-notifications/ping.py API_BASE_URL="https://api.ummc.cosmos.care/v1" BOT_TOKEN= API_AUTHORIZATION_TOKEN= $SERVICE_DIR/python $SERVICE_DIR/ping.py --session 2
